# claude-daily-post.timer
[Unit]
Description=Run claude-daily-post every day at 06:30
Requires=claude-daily-post.service

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
